# debian/panelsdcc-connect.service
[Unit]
Description=panelsDCC-connect - RESTful API for DCC command stations
After=network.target

[Service]
Type=simple
User=dcc-io
Group=dcc-io
WorkingDirectory=/var/lib/panelsdcc-connect
Environment="JMRI_HOME=/opt/jmri/current"
Environment="JAVA_OPTS=-Xmx512m"

# Build classpath with JMRI
# Note: Wildcards in classpath need to be expanded by shell
# JAR_FILENAME_PLACEHOLDER will be replaced during package build
# Fallback: if placeholder not replaced, use wildcard to find JAR
ExecStart=/bin/bash -c 'JAR_FILE="JAR_FILENAME_PLACEHOLDER"; if [ "$$JAR_FILE" = "JAR_FILENAME_PLACEHOLDER" ]; then JAR_FILE=$$(ls /usr/lib/panelsdcc-connect/panelsDCC-connect-*-jar-with-dependencies.jar 2>/dev/null | head -1); fi; CLASSPATH="$$JAR_FILE:/opt/jmri/current/jmri.jar"; if [ -d /opt/jmri/current/lib ]; then for jar in /opt/jmri/current/lib/*.jar; do [ -f "$$jar" ] && CLASSPATH="$$CLASSPATH:$$jar"; done; fi; exec /usr/bin/java -cp "$$CLASSPATH" cc.panelsd.connect.daemon.DccIoDaemon 9000'

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=panelsdcc-connect

# Security settings
# Must be false so the daemon can run `sudo -n ... panelsdcc-connect-updater install` for web UI updates.
# (NoNewPrivileges=true blocks sudo/setuid and causes: "The no new privileges flag is set".)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/panelsdcc-connect

# Allow access to USB, serial devices, and hardware interfaces (GPIO, I2C, SPI)
# netdev: required for nmcli / NetworkManager (WiFi scan, connect, hotspot) as this user;
# systemd uses only this list, not /etc/group, when SupplementaryGroups is set.
SupplementaryGroups=dialout plugdev gpio i2c spi input netdev

# Allow access to device files (USB, serial ports)
PrivateDevices=false

[Install]
WantedBy=multi-user.target

