# debian/prerm
#!/bin/bash
set -e

# Pre-removal script for panelsdcc-connect

# Stop and disable the service before removal
if systemctl is-active panelsdcc-connect >/dev/null 2>&1; then
    echo "Stopping panelsdcc-connect service..."
    systemctl stop panelsdcc-connect
fi

if systemctl is-enabled panelsdcc-connect >/dev/null 2>&1; then
    echo "Disabling panelsdcc-connect service..."
    systemctl disable panelsdcc-connect
fi

# DEBHELPER
exit 0

